# templates/default/rngd.service.erb
#
#  Dynamically generated by Chef on <%= node["fqdn"] %>
#  Local modifications will be overwritten by Chef.
#
[Unit]
Description=Hardware RNG Entropy Gatherer Daemon

[Service]
EnvironmentFile=/etc/sysconfig/rngd
ExecStart=/sbin/rngd -f $EXTRAOPTIONS

[Install]
WantedBy=multi-user.target